# docker/cron/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord_run.log
logfile_maxbytes=50MB
logfile_backups=10

[program:cron]
command=/startup.sh
autostart=true
autorestart=unexpected
exitcodes=0
